# sv/swagger/run
#!/bin/sh

#sv start socklog-unix || exit 1

cd /swagger-editor

exec 2>&1
exec grunt connect:dist
